# etc/systemd-vaultd.socket
[Unit]
Description=systemd-vaultd socket

[Socket]
ListenStream=/run/systemd-vaultd/sock

[Install]
WantedBy=sockets.target
